# systemd/pi-nvr.service
[Unit]
Description=Pi-NVR - lightweight CCTV network video recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__INSTALL_PREFIX__
EnvironmentFile=__ENV_FILE__
ExecStart=__INSTALL_PREFIX__/venv/bin/python3 -m app.main
Restart=on-failure
RestartSec=5

# Hardening. Not exhaustive, but meaningfully reduces blast radius given
# this process shells out to ffmpeg and touches removable storage.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# The install prefix and config dir must exist -- if they don't, something
# is genuinely wrong and the service should fail loudly. The removable-
# media mount points are prefixed with "-" because they're optional: on a
# fresh Raspberry Pi OS install, /media, /mnt, and especially /run/media
# often don't exist yet (no USB drive has ever been plugged in), and
# systemd's ReadWritePaths fails the entire sandbox setup (exit code
# 226/NAMESPACE) if a listed path is missing and not marked optional.
ReadWritePaths=__INSTALL_PREFIX__ /etc/pi-nvr
ReadWritePaths=-/media
ReadWritePaths=-/mnt
ReadWritePaths=-/run/media
ProtectHome=true

# ffmpeg subprocesses + a handful of concurrent camera streams can open
# many file descriptors; raise the default limit.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
